# dhclient-renew.service
[Unit]
Description=Renew DHCP lease for wwan0

[Service]
Type=oneshot
ExecStart=/sbin/dhclient -v wwan0

