# dhcpcd-renew.timer
[Unit]
Description=Renew wwan0 DHCP lease every 15 minutes

[Timer]
OnCalendar=*-*-* *:0/15:00
Persistent=true

[Install]
WantedBy=timers.target

