# dhcpcd-renew.service
[Unit]
Description=Check and release/renew DHCP lease for wwan0

[Service]
Type=oneshot
ExecStart=/home/pi/network_restart.sh

