# dhcpcd-renew.service
[Unit]
Description=Renew DHCP lease for wwan0

[Service]
Type=oneshot
ExecStart=dhcpcd -n wwan0

