# dhcpcd-renew.service
[Unit]
Description=Renew DHCP lease for wwan0

[Service]
Type=oneshot
ExecStart=dhclient -n wwan0

